# docker/all-in-one/resource/supervisord.conf

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket


[program:doris_fe]
command=bash /opt/apache-doris/fe/bin/start_fe.sh
numprocs=1
directory=/opt/apache-doris/fe
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=15
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/opt/apache-doris/fe/log/fe.out
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10


[program:doris_be]
command=bash /opt/apache-doris/be/bin/start_be.sh
numprocs=1
directory=/opt/apache-doris/be
autostart=true
autorestart=true
startsecs=5
startretries=3
# do force kill before BE is good with graceful shutdown
stopsignal=KILL
stopwaitsecs=15
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/opt/apache-doris/be/log/be.out
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10


[program:deploy]
command=bash /opt/deploy.sh
numprocs=1
directory=/opt
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=15
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
